# docker/util/ld.so.preload
# this file enables opengl with nvidia gpu.
/usr/local/$LIB/libGL.so.1
/usr/local/$LIB/libEGL.so.1
